# systemd/hermes-mtf-macd-tuner.timer
[Unit]
Description=Hermes MTF-MACD Tuner — quick re-test every 1min, full sweep every 24h
After=network.target

[Timer]
OnBootSec=5
OnUnitActiveSec=1m
Persistent=true

[Install]
WantedBy=timers.target
