# systemd/hermes-ma-cross-5m-tuner.timer
[Unit]
Description=Hermes MA Cross 5M Tuner — full sweep every 24h
After=network.target

[Timer]
OnBootSec=5
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
